# etc/profile.d/qt3-devel.csh
setenv QTDIR /usr/lib/qt3
setenv QMAKESPEC /usr/lib/qt3/mkspecs/cygwin-g++
