# etc/profile.d/qt3-devel.sh
export QTDIR=/usr/lib/qt3
export QMAKESPEC=/usr/lib/qt3/mkspecs/cygwin-g++
